# live-build/config/includes.chroot/etc/skel/.config/autostart/nexus-installer.desktop
[Desktop Entry]
Type=Application
Name=Install Nexus Linux
Comment=Launches the Nexus Linux installer automatically in the live session
Exec=/usr/local/bin/nexus-installer
Icon=system-software-install
Terminal=false
X-GNOME-Autostart-enabled=true
NoDisplay=false
